# service/hostapd/run
#!/bin/sh

exec hostapd -d hostapd.conf
